# crontab
# NetroHub Cron Jobs
# Run Laravel scheduler every minute
* * * * * cd /path/to/netrohub && php artisan schedule:run >> /dev/null 2>&1

# Clean up old log files (run daily at 2 AM)
0 2 * * * find /path/to/netrohub/storage/logs -name "*.log" -mtime +30 -delete

# Clean up old temporary files (run daily at 3 AM)
0 3 * * * find /path/to/netrohub/storage/app/temp -type f -mtime +7 -delete

# Generate sitemap (run daily at 4 AM)
0 4 * * * cd /path/to/netrohub && php artisan sitemap:generate >> /dev/null 2>&1
